# php7-fpm/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/php/
-rw-r--r--	root/root	etc/php/php-fpm.conf
-rw-r--r--	root/root	etc/php/php-fpm.conf.default
drwxr-xr-x	root/root	etc/php/php-fpm.d/
-rw-r--r--	root/root	etc/php/php-fpm.d/www.conf.default
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr-x	root/root	etc/rc.d/php-fpm
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/php/
drwxr-xr-x	root/root	usr/php/man/
drwxr-xr-x	root/root	usr/php/man/man8/
-rw-r--r--	root/root	usr/php/man/man8/php-fpm.8.gz
drwxr-xr-x	root/root	usr/php/php/
drwxr-xr-x	root/root	usr/php/php/fpm/
-rw-r--r--	root/root	usr/php/php/fpm/status.html
drwxr-xr-x	root/root	usr/sbin/
-rwxr-xr-x	root/root	usr/sbin/php-fpm
drwxr-xr-x	root/root	usr/share/
drwxr-xr-x	root/root	usr/share/man/
drwxr-xr-x	root/root	usr/share/man/man8/
-rw-r--r--	root/root	usr/share/man/man8/php-fpm.8.gz
drwxr-xr-x	root/root	var/
drwxr-xr-x	root/root	var/log/
drwxr-xr-x	root/root	var/run/
