# tlp/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/init.d/
-rwxr-xr-x	root/root	etc/init.d/tlp
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr-x	root/root	etc/rc.d/tlp
drwxr-xr-x	root/root	lib/
drwxr-xr-x	root/root	lib/udev/
drwxr-xr-x	root/root	lib/udev/rules.d/
-rw-r--r--	root/root	lib/udev/rules.d/85-tlp.rules
-rwxr-xr-x	root/root	lib/udev/tlp-usb-udev
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/bin/
-rwxr-xr-x	root/root	usr/bin/bluetooth
-rwxr-xr-x	root/root	usr/bin/run-on-ac
lrwxrwxrwx	root/root	usr/bin/run-on-bat -> run-on-ac
-rwxr-xr-x	root/root	usr/bin/tlp
-rwxr-xr-x	root/root	usr/bin/tlp-pcilist
-rwxr-xr-x	root/root	usr/bin/tlp-stat
-rwxr-xr-x	root/root	usr/bin/tlp-usblist
lrwxrwxrwx	root/root	usr/bin/wifi -> bluetooth
lrwxrwxrwx	root/root	usr/bin/wwan -> bluetooth
drwxr-xr-x	root/root	usr/etc/
-rw-r--r--	root/root	usr/etc/tlp
drwxr-xr-x	root/root	usr/share/
drwxr-xr-x	root/root	usr/share/bash-completion/
drwxr-xr-x	root/root	usr/share/bash-completion/completions/
lrwxrwxrwx	root/root	usr/share/bash-completion/completions/bluetooth -> tlp
-rw-r--r--	root/root	usr/share/bash-completion/completions/tlp
lrwxrwxrwx	root/root	usr/share/bash-completion/completions/tlp-stat -> tlp
lrwxrwxrwx	root/root	usr/share/bash-completion/completions/wifi -> tlp
lrwxrwxrwx	root/root	usr/share/bash-completion/completions/wwan -> tlp
drwxr-xr-x	root/root	usr/share/man/
drwxr-xr-x	root/root	usr/share/man/man1/
-rw-r--r--	root/root	usr/share/man/man1/bluetooth.1.gz
-rw-r--r--	root/root	usr/share/man/man1/run-on-ac.1.gz
-rw-r--r--	root/root	usr/share/man/man1/run-on-bat.1.gz
-rw-r--r--	root/root	usr/share/man/man1/tlp-pcilist.1.gz
-rw-r--r--	root/root	usr/share/man/man1/tlp-usblist.1.gz
-rw-r--r--	root/root	usr/share/man/man1/wifi.1.gz
-rw-r--r--	root/root	usr/share/man/man1/wwan.1.gz
drwxr-xr-x	root/root	usr/share/man/man8/
-rw-r--r--	root/root	usr/share/man/man8/tlp-sleep.service.8.gz
-rw-r--r--	root/root	usr/share/man/man8/tlp-stat.8.gz
-rw-r--r--	root/root	usr/share/man/man8/tlp.8.gz
-rw-r--r--	root/root	usr/share/man/man8/tlp.service.8.gz
drwxr-xr-x	root/root	usr/share/metainfo/
-rw-r--r--	root/root	usr/share/metainfo/de.linrunner.tlp.metainfo.xml
drwxr-xr-x	root/root	usr/share/tlp/
-rwxr-xr-x	root/root	usr/share/tlp/tlp-functions
-rwxr-xr-x	root/root	usr/share/tlp/tlp-rf-func
-rwxr-xr-x	root/root	usr/share/tlp/tpacpi-bat
