# bluez5/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/dbus-1/
drwxr-xr-x	root/root	etc/dbus-1/system.d/
-rw-r--r--	root/root	etc/dbus-1/system.d/bluetooth.conf
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr-x	root/root	etc/rc.d/bluetoothd
drwxr-xr-x	root/root	lib/
drwxr-xr-x	root/root	lib/udev/
-rwxr-xr-x	root/root	lib/udev/hid2hci
drwxr-xr-x	root/root	lib/udev/rules.d/
-rw-r--r--	root/root	lib/udev/rules.d/97-hid2hci.rules
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/bin/
-rwxr-xr-x	root/root	usr/bin/bccmd
-rwxr-xr-x	root/root	usr/bin/bluemoon
-rwxr-xr-x	root/root	usr/bin/bluetoothctl
-rwxr-xr-x	root/root	usr/bin/btattach
-rwxr-xr-x	root/root	usr/bin/btmon
-rwxr-xr-x	root/root	usr/bin/hex2hcd
-rwxr-xr-x	root/root	usr/bin/l2ping
-rwxr-xr-x	root/root	usr/bin/l2test
-rwxr-xr-x	root/root	usr/bin/mpris-proxy
-rwxr-xr-x	root/root	usr/bin/rctest
drwxr-xr-x	root/root	usr/include/
drwxr-xr-x	root/root	usr/include/bluetooth/
-rw-r--r--	root/root	usr/include/bluetooth/bluetooth.h
-rw-r--r--	root/root	usr/include/bluetooth/bnep.h
-rw-r--r--	root/root	usr/include/bluetooth/cmtp.h
-rw-r--r--	root/root	usr/include/bluetooth/hci.h
-rw-r--r--	root/root	usr/include/bluetooth/hci_lib.h
-rw-r--r--	root/root	usr/include/bluetooth/hidp.h
-rw-r--r--	root/root	usr/include/bluetooth/l2cap.h
-rw-r--r--	root/root	usr/include/bluetooth/rfcomm.h
-rw-r--r--	root/root	usr/include/bluetooth/sco.h
-rw-r--r--	root/root	usr/include/bluetooth/sdp.h
-rw-r--r--	root/root	usr/include/bluetooth/sdp_lib.h
drwxr-xr-x	root/root	usr/lib/
drwxr-xr-x	root/root	usr/lib/bluez/
drwxr-xr-x	root/root	usr/lib/bluez/test/
-rwxr-xr-x	root/root	usr/lib/bluez/test/bluezutils.py
-rwxr-xr-x	root/root	usr/lib/bluez/test/dbusdef.py
-rwxr-xr-x	root/root	usr/lib/bluez/test/example-advertisement
-rwxr-xr-x	root/root	usr/lib/bluez/test/example-gatt-client
-rwxr-xr-x	root/root	usr/lib/bluez/test/example-gatt-server
-rwxr-xr-x	root/root	usr/lib/bluez/test/ftp-client
-rwxr-xr-x	root/root	usr/lib/bluez/test/list-devices
-rwxr-xr-x	root/root	usr/lib/bluez/test/map-client
-rwxr-xr-x	root/root	usr/lib/bluez/test/monitor-bluetooth
-rwxr-xr-x	root/root	usr/lib/bluez/test/opp-client
-rwxr-xr-x	root/root	usr/lib/bluez/test/pbap-client
-rwxr-xr-x	root/root	usr/lib/bluez/test/sap_client.py
-rwxr-xr-x	root/root	usr/lib/bluez/test/service-did.xml
-rwxr-xr-x	root/root	usr/lib/bluez/test/service-ftp.xml
-rwxr-xr-x	root/root	usr/lib/bluez/test/service-opp.xml
-rwxr-xr-x	root/root	usr/lib/bluez/test/service-record.dtd
-rwxr-xr-x	root/root	usr/lib/bluez/test/service-spp.xml
-rwxr-xr-x	root/root	usr/lib/bluez/test/simple-agent
-rwxr-xr-x	root/root	usr/lib/bluez/test/simple-endpoint
-rwxr-xr-x	root/root	usr/lib/bluez/test/simple-player
-rwxr-xr-x	root/root	usr/lib/bluez/test/test-adapter
-rwxr-xr-x	root/root	usr/lib/bluez/test/test-device
-rwxr-xr-x	root/root	usr/lib/bluez/test/test-discovery
-rwxr-xr-x	root/root	usr/lib/bluez/test/test-gatt-profile
-rwxr-xr-x	root/root	usr/lib/bluez/test/test-health
-rwxr-xr-x	root/root	usr/lib/bluez/test/test-health-sink
-rwxr-xr-x	root/root	usr/lib/bluez/test/test-hfp
-rwxr-xr-x	root/root	usr/lib/bluez/test/test-manager
-rwxr-xr-x	root/root	usr/lib/bluez/test/test-nap
-rwxr-xr-x	root/root	usr/lib/bluez/test/test-network
-rwxr-xr-x	root/root	usr/lib/bluez/test/test-profile
-rwxr-xr-x	root/root	usr/lib/bluez/test/test-sap-server
drwxr-xr-x	root/root	usr/lib/bluez5/
drwxr-xr-x	root/root	usr/lib/bluez5/bluetooth/
-rwxr-xr-x	root/root	usr/lib/bluez5/bluetooth/bluetoothd
-rwxr-xr-x	root/root	usr/lib/bluez5/bluetooth/obexd
drwxr-xr-x	root/root	usr/lib/cups/
drwxr-xr-x	root/root	usr/lib/cups/backend/
-rwxr-xr-x	root/root	usr/lib/cups/backend/bluetooth
-rwxr-xr-x	root/root	usr/lib/libbluetooth.la
lrwxrwxrwx	root/root	usr/lib/libbluetooth.so -> libbluetooth.so.3.18.16
lrwxrwxrwx	root/root	usr/lib/libbluetooth.so.3 -> libbluetooth.so.3.18.16
-rwxr-xr-x	root/root	usr/lib/libbluetooth.so.3.18.16
drwxr-xr-x	root/root	usr/lib/pkgconfig/
-rw-r--r--	root/root	usr/lib/pkgconfig/bluez.pc
drwxr-xr-x	root/root	usr/share/
drwxr-xr-x	root/root	usr/share/man/
drwxr-xr-x	root/root	usr/share/man/man1/
-rw-r--r--	root/root	usr/share/man/man1/bccmd.1.gz
-rw-r--r--	root/root	usr/share/man/man1/btattach.1.gz
-rw-r--r--	root/root	usr/share/man/man1/hid2hci.1.gz
-rw-r--r--	root/root	usr/share/man/man1/l2ping.1.gz
-rw-r--r--	root/root	usr/share/man/man1/rctest.1.gz
drwxr-xr-x	root/root	usr/share/man/man8/
-rw-r--r--	root/root	usr/share/man/man8/bluetoothd.8.gz
