# distcc/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr-x	root/root	etc/rc.d/distccd
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/bin/
-rwxr-xr-x	root/root	usr/bin/distcc
-rwxr-xr-x	root/root	usr/bin/distccmon-text
-rwxr-xr-x	root/root	usr/bin/lsdistcc
-rwxr-xr-x	root/root	usr/bin/pump
drwxr-xr-x	root/root	usr/lib/
drwxr-xr-x	root/root	usr/lib/distcc/
lrwxrwxrwx	root/root	usr/lib/distcc/c++ -> ../../bin/distcc
lrwxrwxrwx	root/root	usr/lib/distcc/cc -> ../../bin/distcc
lrwxrwxrwx	root/root	usr/lib/distcc/cpp -> ../../bin/distcc
lrwxrwxrwx	root/root	usr/lib/distcc/g++ -> ../../bin/distcc
lrwxrwxrwx	root/root	usr/lib/distcc/gcc -> ../../bin/distcc
drwxr-xr-x	root/root	usr/lib/python3.6/
drwxr-xr-x	root/root	usr/lib/python3.6/site-packages/
-rw-r--r--	root/root	usr/lib/python3.6/site-packages/include_server-3.3.2-py3.6.egg-info
drwxr-xr-x	root/root	usr/lib/python3.6/site-packages/include_server/
drwxr-xr-x	root/root	usr/lib/python3.6/site-packages/include_server/__pycache__/
-rw-r--r--	root/root	usr/lib/python3.6/site-packages/include_server/__pycache__/basics.cpython-36.pyc
-rw-r--r--	root/root	usr/lib/python3.6/site-packages/include_server/__pycache__/basics_test.cpython-36.pyc
-rw-r--r--	root/root	usr/lib/python3.6/site-packages/include_server/__pycache__/c_extensions_test.cpython-36.pyc
-rw-r--r--	root/root	usr/lib/python3.6/site-packages/include_server/__pycache__/cache_basics.cpython-36.pyc
-rw-r--r--	root/root	usr/lib/python3.6/site-packages/include_server/__pycache__/compiler_defaults.cpython-36.pyc
-rw-r--r--	root/root	usr/lib/python3.6/site-packages/include_server/__pycache__/compress_files.cpython-36.pyc
-rw-r--r--	root/root	usr/lib/python3.6/site-packages/include_server/__pycache__/include_analyzer.cpython-36.pyc
-rw-r--r--	root/root	usr/lib/python3.6/site-packages/include_server/__pycache__/include_analyzer_memoizing_node.cpython-36.pyc
-rw-r--r--	root/root	usr/lib/python3.6/site-packages/include_server/__pycache__/include_analyzer_memoizing_node_test.cpython-36.pyc
-rw-r--r--	root/root	usr/lib/python3.6/site-packages/include_server/__pycache__/include_analyzer_test.cpython-36.pyc
-rw-r--r--	root/root	usr/lib/python3.6/site-packages/include_server/__pycache__/include_server.cpython-36.pyc
-rw-r--r--	root/root	usr/lib/python3.6/site-packages/include_server/__pycache__/include_server_test.cpython-36.pyc
-rw-r--r--	root/root	usr/lib/python3.6/site-packages/include_server/__pycache__/macro_eval.cpython-36.pyc
-rw-r--r--	root/root	usr/lib/python3.6/site-packages/include_server/__pycache__/macro_eval_test.cpython-36.pyc
-rw-r--r--	root/root	usr/lib/python3.6/site-packages/include_server/__pycache__/mirror_path.cpython-36.pyc
-rw-r--r--	root/root	usr/lib/python3.6/site-packages/include_server/__pycache__/mirror_path_test.cpython-36.pyc
-rw-r--r--	root/root	usr/lib/python3.6/site-packages/include_server/__pycache__/parse_command.cpython-36.pyc
-rw-r--r--	root/root	usr/lib/python3.6/site-packages/include_server/__pycache__/parse_command_test.cpython-36.pyc
-rw-r--r--	root/root	usr/lib/python3.6/site-packages/include_server/__pycache__/parse_file.cpython-36.pyc
-rw-r--r--	root/root	usr/lib/python3.6/site-packages/include_server/__pycache__/parse_file_test.cpython-36.pyc
-rw-r--r--	root/root	usr/lib/python3.6/site-packages/include_server/__pycache__/run.cpython-36.pyc
-rw-r--r--	root/root	usr/lib/python3.6/site-packages/include_server/__pycache__/setup.cpython-36.pyc
-rw-r--r--	root/root	usr/lib/python3.6/site-packages/include_server/__pycache__/statistics.cpython-36.pyc
-rw-r--r--	root/root	usr/lib/python3.6/site-packages/include_server/basics.py
-rw-r--r--	root/root	usr/lib/python3.6/site-packages/include_server/basics_test.py
-rw-r--r--	root/root	usr/lib/python3.6/site-packages/include_server/c_extensions_test.py
-rw-r--r--	root/root	usr/lib/python3.6/site-packages/include_server/cache_basics.py
-rw-r--r--	root/root	usr/lib/python3.6/site-packages/include_server/compiler_defaults.py
-rw-r--r--	root/root	usr/lib/python3.6/site-packages/include_server/compress_files.py
-rwxr-xr-x	root/root	usr/lib/python3.6/site-packages/include_server/distcc_pump_c_extensions.cpython-36m-x86_64-linux-gnu.so
-rw-r--r--	root/root	usr/lib/python3.6/site-packages/include_server/include_analyzer.py
-rw-r--r--	root/root	usr/lib/python3.6/site-packages/include_server/include_analyzer_memoizing_node.py
-rw-r--r--	root/root	usr/lib/python3.6/site-packages/include_server/include_analyzer_memoizing_node_test.py
-rw-r--r--	root/root	usr/lib/python3.6/site-packages/include_server/include_analyzer_test.py
-rw-r--r--	root/root	usr/lib/python3.6/site-packages/include_server/include_server.py
-rw-r--r--	root/root	usr/lib/python3.6/site-packages/include_server/include_server_test.py
-rw-r--r--	root/root	usr/lib/python3.6/site-packages/include_server/macro_eval.py
-rw-r--r--	root/root	usr/lib/python3.6/site-packages/include_server/macro_eval_test.py
-rw-r--r--	root/root	usr/lib/python3.6/site-packages/include_server/mirror_path.py
-rw-r--r--	root/root	usr/lib/python3.6/site-packages/include_server/mirror_path_test.py
-rw-r--r--	root/root	usr/lib/python3.6/site-packages/include_server/parse_command.py
-rw-r--r--	root/root	usr/lib/python3.6/site-packages/include_server/parse_command_test.py
-rw-r--r--	root/root	usr/lib/python3.6/site-packages/include_server/parse_file.py
-rw-r--r--	root/root	usr/lib/python3.6/site-packages/include_server/parse_file_test.py
-rw-r--r--	root/root	usr/lib/python3.6/site-packages/include_server/run.py
-rw-r--r--	root/root	usr/lib/python3.6/site-packages/include_server/setup.py
-rw-r--r--	root/root	usr/lib/python3.6/site-packages/include_server/statistics.py
drwxr-xr-x	root/root	usr/sbin/
-rwxr-xr-x	root/root	usr/sbin/distccd
-rwxr-xr-x	root/root	usr/sbin/update-distcc-symlinks
drwxr-xr-x	root/root	usr/share/
drwxr-xr-x	root/root	usr/share/man/
drwxr-xr-x	root/root	usr/share/man/man1/
-rw-r--r--	root/root	usr/share/man/man1/distcc.1.gz
-rw-r--r--	root/root	usr/share/man/man1/distccd.1.gz
-rw-r--r--	root/root	usr/share/man/man1/distccmon-text.1.gz
-rw-r--r--	root/root	usr/share/man/man1/include_server.1.gz
-rw-r--r--	root/root	usr/share/man/man1/lsdistcc.1.gz
-rw-r--r--	root/root	usr/share/man/man1/pump.1.gz
